# systemd/user/arch-update-tray.service.d/path-override.conf
[Service]
Environment=PATH=/usr/bin:/bin
