# systemd/user/steam-autopicture.service
[Unit]
Description=Auto Big Picture — Abre Steam Big Picture al conectar un mando
After=graphical-session.target
BindsTo=graphical-session.target

[Service]
Type=simple
ExecStart=%h/scripts/steam_autopicture.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
